# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
command=uvicorn src.urbanpulse.api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api.log

[program:dashboard]
command=streamlit run app.py --server.port=7860 --server.address=0.0.0.0 --server.headless=true
directory=/app
autostart=true
autorestart=true
environment=API_BASE_URL="http://localhost:8000/api/v1"
stdout_logfile=/tmp/dashboard.log
stderr_logfile=/tmp/dashboard.log
