# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:api]
command=uvicorn urbanpulse.api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dashboard]
command=streamlit run src/urbanpulse/dashboard/app.py --server.port=7860 --server.address=0.0.0.0 --server.headless=true
directory=/app
autostart=true
autorestart=true
environment=API_BASE_URL="http://localhost:8000/api/v1",DATABASE_URL="%(ENV_DATABASE_URL)s",OPENAQ_API_KEY="%(ENV_OPENAQ_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
